# client-build/resources/linux/gexd_centos.conf
[program:gexd]
command=/usr/lib/gex/java/bin/java -jar /usr/lib/gex/gexd.jar
user=gexuser
autostart=true
autorestart=true
stderr_logfile=/var/log/gex.err.log
stdout_logfile=/var/log/gex.out.log

